# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=python -m avatarfactory.daemon_runner --mode full --host 127.0.0.1 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:monitor]
command=node /app/monitor/server/entry.mjs
directory=/app/monitor
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="4321",API_BASE_URL="http://127.0.0.1:8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:journal]
command=node /app/journal/server/entry.mjs
directory=/app/journal
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="4322",API_BASE_URL="http://127.0.0.1:8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=35

[program:admin]
command=node /app/admin/server/entry.mjs
directory=/app/admin
autostart=true
autorestart=true
environment=HOST="127.0.0.1",PORT="4323",API_BASE_URL="http://127.0.0.1:8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40

[program:nginx]
command=nginx -c /app/scripts/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
